# sys-apps/world/files/world.timer
[Unit]
Description=Daily update of overlays

[Timer]
OnStartupSec=1min
OnUnitActiveSec=1d

[Install]
WantedBy=timers.target
