# sys-apps/world/files/world.service
[Unit]
Description=service to generate repositories list

[Service]
ExecStart=r

[Install]
WantedBy=multi-user.target
